# example/crontab
# min hour dom month dow command
  0   *    *   *     *   PATH=$HOME/bin:$HOME/.cargo/bin:$PATH; exec update-feeds ~/.feeds
